# systemd/.local/share/systemd/user/stand.timer
[Unit]
Description="Timer to remind me to change position"

[Timer]
OnUnitActiveSec=3600
OnBootSec=60

[Install]
WantedBy=timers.target
